# dist/rgb-controller.service
[Unit]
Description=RGB Controller Ambilight

[Service]
Type=simple
ExecStart=%h/.local/bin/rgb-controller-launcher
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
